# meta-refkit/recipes-security/nftables-settings-default/files/firewall.path
[Unit]
Description=Track changes to firewall configuration

[Path]
PathChanged=/run/firewall/firewall.ruleset

[Install]
WantedBy=network.target
